# files/consul.conf
description "Consul server process"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

setuid consul
setgid consul

exec /opt/consul/consul agent -config-file /etc/consul.d/consul.json
