# files/nomad.conf
description "Nomad server process"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

setuid nomad
setgid nomad

exec /opt/nomad/nomad agent -config /etc/nomad.d/nomad.json
